# daemon/systemd/double-ctrl.service
[Unit]
Description=Double Ctrl Hotkey Daemon (User Session)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/double-ctrl
Restart=on-failure
# Set log level (info=default, debug=verbose, warn=minimal, error=errors only)
Environment="RUST_LOG=info"

[Install]
WantedBy=default.target
